# Library/LaunchAgents/com.bratteng.updateDailyExtension.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.bratteng.updateDailyExtension</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/ombratteng/.local/bin/updateDailyExtension</string>
	</array>

	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Hour</key>
			<integer>10</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Hour</key>
			<integer>18</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>StandardErrorPath</key>
	<string>/tmp/com.bratteng.updateDailyExtension.err</string>

	<key>StandardOutPath</key>
	<string>/tmp/com.bratteng.updateDailyExtension.out</string>
</dict>
</plist>
